# debian/powa-collector.service
[Unit]
Description=PoWA Collector daemon for remote snapshots

[Service]
User=powa
ExecStart=/usr/bin/powa-collector
ExecReload=kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
